# personal_budget/horizon-api/deploy/horizon-backup.timer
[Unit]
Description=Horizon — run DB backup daily at 03:30

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=horizon-backup.service

[Install]
WantedBy=timers.target
